# ops/com.jonathangu.ocbrain.autopilot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.jonathangu.ocbrain.autopilot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/uv</string>
    <string>run</string>
    <string>--directory</string>
    <string>/Users/guclaw/.openclaw/workspace/ocbrain</string>
    <string>ocbrain</string>
    <string>autopilot</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>OCBRAIN_DB</key>
    <string>/Users/guclaw/.openclaw/workspace/ocbrain/data/ocbrain.sqlite</string>
  </dict>
  <key>StartInterval</key><integer>1800</integer>
  <key>RunAtLoad</key><false/>
  <key>ProcessType</key><string>Background</string>
  <key>LowPriorityBackgroundIO</key><true/>
  <key>Nice</key><integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/guclaw/.openclaw/workspace/ocbrain/data/logs/autopilot.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/guclaw/.openclaw/workspace/ocbrain/data/logs/autopilot.err.log</string>
</dict>
</plist>
