# ops/com.jonathangu.ocbrain.stallcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.jonathangu.ocbrain.stallcheck</string>
  <!--
    DEPRECATED / RETIRED — DO NOT INSTALL OR LOAD.

    This label is retained only so operators can identify and remove an old
    installation. v0.4.1 has no recurring stallcheck or pager job. Diagnostic
    code remains available only for an explicit manual invocation. There is
    deliberately no timer and the inert program cannot dispatch work.
  -->
  <key>Disabled</key><true/>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/false</string>
  </array>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
